# deploy/ansible/roles/deploy/templates/pombot.conf
description "daemon for pombot app"

start on startup
stop on shutdown
respawn

script
  {{which_nodejs.stdout | default('/usr/bin/nodejs')}} {{site_path}}/app.js
end script
